# packaging/debian/postrm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1 && [ "$1" = "remove" ]; then
    if command -v loginctl >/dev/null 2>&1; then
        loginctl list-users --no-legend 2>/dev/null | awk '{print $2}' | while read -r user; do
            [ -z "$user" ] && continue
            [ "$user" = "root" ] && continue
            runuser -l "$user" -c "systemctl --user disable armonix.service armonix-gui.service >/dev/null 2>&1 || true" || true
            runuser -l "$user" -c "systemctl --user daemon-reload >/dev/null 2>&1 || true" || true
        done || true
    fi
fi

exit 0
