# packaging/debian/postinst
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1 && [ "$1" = "configure" ]; then
    # Enable and start the user service for every logged-in non-root user.
    # On a dedicated single-user machine this is typically just one user.
    if command -v loginctl >/dev/null 2>&1; then
        loginctl list-users --no-legend 2>/dev/null | awk '{print $2}' | while read -r user; do
            [ -z "$user" ] && continue
            [ "$user" = "root" ] && continue
            runuser -l "$user" -c "systemctl --user daemon-reload >/dev/null 2>&1 || true" || true
            runuser -l "$user" -c "systemctl --user enable --now armonix-gui.service >/dev/null 2>&1 || true" || true
        done || true
    fi
fi

exit 0
