# packaging/armonix-gui.service
[Unit]
Description=Armonix MIDI controller with GUI (user service)
After=graphical-session.target sound.target

[Service]
Type=simple
ExecStart=/usr/bin/armonix-gui
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
