# packaging/armonix.service
[Unit]
Description=Armonix MIDI controller
After=network.target sound.target

[Service]
Type=simple
User=armonix
Group=armonix
SupplementaryGroups=input
WorkingDirectory=/usr/lib/armonix
ExecStart=/usr/bin/armonix-engine
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
